# systemd/meiso-edge.service
[Unit]
Description=Meiso Edge Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/meiso_glass
ExecStart=/opt/meiso_glass/.venv/bin/meiso edge --config /opt/meiso_glass/configs/examples/edge.generic.yaml
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
